# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:phoenix]
command=phoenix serve
directory=%(here)s
autostart=true
autorestart=true
stderr_logfile=/tmp/phoenix_err.log
stdout_logfile=/tmp/phoenix_out.log
user=%(ENV_USER)s

[program:main]
command=python main.py
directory=%(here)s
autostart=true
autorestart=true
stderr_logfile=/tmp/main_err.log
stdout_logfile=/tmp/main_out.log
user=%(ENV_USER)s